# skel/etc/cron.d/connection-test
MAILTO=root

*/10 *  * * *   root    /usr/local/sbin/connection-test  >>  /var/log/connection-test.log  2>&1

