# skel/etc/cron.d/sync_clock.cron
MAILTO=root

*/10 *  * * *   root    /etc/sync_clock

