# skel/etc/rc.local
#!/bin/sh

s="`hostname` startup complete `date -u +'%Y/%m/%d %T'`UTC"
echo "${s}" | mail -s  "${s}" root

